# deploy/daybreak.service
[Unit]
Description=Daybreak ORB bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/daybreak
EnvironmentFile=/opt/daybreak/.env
ExecStart=/opt/daybreak/.venv/bin/python scripts/run.py config/config.yaml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
